# ops/server/komui-traffic-switch.service
[Unit]
Description=Apply KOMUI production traffic switch request
After=nginx.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/komui-traffic-switch-apply
